# debian/jd-price.service
[Unit]
Description=JD Price Comparison Service
After=network.target

[Service]
Type=simple
User=jdprice
ExecStart=/usr/bin/jd-price
Restart=on-failure
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
